# com.local.mlx-native-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.local.mlx-native-server</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jscheel/tools/mlx-tools/start-server.sh</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/jscheel/tools/mlx-tools</string>
    
    <!-- Run at load (start automatically) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep the service alive (restart if it crashes) -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/jscheel/tools/mlx-tools/logs/stdout.log</string>
    
    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/jscheel/tools/mlx-tools/logs/stderr.log</string>
    
    <!-- Environment variables (if needed) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Throttle interval (seconds to wait before restart after crash) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
